# container/services/odoo/run
#!/bin/sh
exec 2>&1

# Directories and init
mkdir -p /var/log/odoo
chown -R odoo. /var/log/odoo

exec chpst -u odoo odoo -c /container/config/odoo/odoo.conf
