# container/services/redis/run
#!/bin/sh
exec 2>&1

exec redis-server /container/config/redis/redis/conf
